# otx-tor.timer
[Unit]
Description=Run otx-tor script

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
